# rootfs/etc/services.d/rtl_433/run
#!/usr/bin/with-contenv bashio
declare confdir
declare custdir
declare cfig
declare flexfig

confdir="/rattler/"
custdir="/config/$(bashio::config 'customize.folder')/"
dir="${confdir}"
if [ $(bashio::config.true 'customize.active') ]; then
  if [ ! -d "${custdir}" ]; then mkdir -p "${custdir}"; fi
  if [ ! -f "${custdir}rtl_433.conf" ]; then cp "${confdir}rtl_433.conf" "${custdir}"; fi
  if [ ! -f "${custdir}output.conf" ]; then cp "${confdir}output.conf" "${custdir}"; fi
  dir="${custdir}"
  if [ $(bashio::config.true 'customize.flex') ]; then
    if [ ! -d "${custdir}flex" ]; then cp "${confdir}flex" "${custdir}" || bashio::exit.nok "Could not create Rattler Flex config folder!"; fi
    for cfig in "${custdir}flex/*.conf"; do
      flexfig="$flexfig -c $cfig "
    done
  fi
fi

rtl_433 -c ${dir}rtl_433.conf $flexfig -c ${dir}output.conf
